# scripts/systemd/ernos-auth-monitor.service
[Unit]
Description=ErnOS Auth Expiry Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/home/admin/ernos/scripts/auth-monitor.sh
# Configure notification channels via environment
Environment=WARN_HOURS=2
# Environment=NOTIFY_PHONE=+1234567890
# Environment=NOTIFY_NTFY=ernos-alerts

[Install]
WantedBy=default.target
